# src/ota_daemon/deploy/ota_daemon.service
[Unit]
Description=OTA Update Daemon (gRPC :50052)
Documentation=file:///opt/robot/ota/README.md
After=network-online.target
Wants=network-online.target

# OTA daemon 独立于导航栈, 不依赖 navigation.service
# 这样它可以停止/重启 navigation.service 来完成 COLD 更新

[Service]
Type=simple
User=sunrise
Group=sunrise
WorkingDirectory=/opt/robot

ExecStart=/opt/robot/bin/ota_daemon --config /opt/robot/ota/ota_daemon.yaml

# 崩溃恢复: 检查未完成的事务
ExecStartPre=/bin/bash -c '\
  for txn in /opt/robot/ota/backup/txn_*.json; do \
    [ -f "$txn" ] || continue; \
    echo "OTA daemon: recovering interrupted transaction: $txn"; \
  done'

Restart=on-failure
RestartSec=3
StartLimitBurst=5
StartLimitIntervalSec=60

# 安全加固
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/robot /home/sunrise/models /home/sunrise/data/SLAM/navigation/install /home/sunrise/data/SLAM/navigation/maps /home/sunrise/data/SLAM/navigation/src/global_planning
ProtectHome=true
PrivateTmp=true

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ota_daemon

[Install]
WantedBy=multi-user.target
