# systemd/nav-planning.service
[Unit]
Description=Navigation - Global Planning (Localizer + PCT Planner)
After=nav-slam.service
Requires=nav-slam.service

[Service]
Type=simple
User=sunrise
Group=sunrise
WorkingDirectory=/opt/nav
ExecStart=/bin/bash /opt/nav/scripts/services/nav-planning.sh
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=15
WatchdogSec=30
KillMode=mixed
KillSignal=SIGINT
SuccessExitStatus=130 SIGINT

# 关键服务保护: 降低 OOM killer 优先级
OOMScoreAdjust=-500
MemoryMax=3G

# 地图和初始位姿可通过 EnvironmentFile 覆盖
EnvironmentFile=-/etc/nav/planning.env
Environment=HOME=/home/sunrise
Environment=LANG=en_US.UTF-8

StandardOutput=journal
StandardError=journal
SyslogIdentifier=nav-planning

[Install]
WantedBy=multi-user.target
