# systemd/nav-lidar.service
[Unit]
Description=Navigation - Livox LiDAR Driver
After=nav-lidar-network.service
Requires=nav-lidar-network.service

[Service]
Type=simple
User=sunrise
Group=sunrise
WorkingDirectory=/home/sunrise/data/SLAM/navigation
ExecStart=/bin/bash /home/sunrise/data/SLAM/navigation/scripts/services/nav-lidar.sh
Restart=on-failure
RestartSec=3
TimeoutStopSec=10
KillMode=mixed
KillSignal=SIGINT
SuccessExitStatus=130 SIGINT

# 环境
Environment=HOME=/home/sunrise
Environment=LANG=en_US.UTF-8

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nav-lidar

[Install]
WantedBy=multi-user.target
