# systemd/nav-planning.service
[Unit]
Description=Navigation - Global Planning (Localizer + PCT Planner)
After=nav-slam.service
Requires=nav-slam.service

[Service]
Type=simple
User=sunrise
Group=sunrise
WorkingDirectory=/opt/nav
ExecStart=/bin/bash /opt/nav/scripts/services/nav-planning.sh
Restart=on-failure
RestartSec=5
TimeoutStopSec=15
KillMode=mixed
KillSignal=SIGINT
SuccessExitStatus=130 SIGINT

# 地图和初始位姿可通过 EnvironmentFile 覆盖
EnvironmentFile=-/etc/nav/planning.env
Environment=HOME=/home/sunrise
Environment=LANG=en_US.UTF-8

StandardOutput=journal
StandardError=journal
SyslogIdentifier=nav-planning

[Install]
WantedBy=multi-user.target
