# systemd/nav-slam.service
[Unit]
Description=Navigation - Fast-LIO2 SLAM
After=nav-lidar.service
Requires=nav-lidar.service

[Service]
Type=simple
User=sunrise
Group=sunrise
WorkingDirectory=/opt/nav
ExecStart=/bin/bash /opt/nav/scripts/services/nav-slam.sh
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=10
WatchdogSec=30
KillMode=mixed
KillSignal=SIGINT
SuccessExitStatus=130 SIGINT

# 关键服务保护: 降低 OOM killer 优先级
OOMScoreAdjust=-500
MemoryMax=4G

Environment=HOME=/home/sunrise
Environment=LANG=en_US.UTF-8

StandardOutput=journal
StandardError=journal
SyslogIdentifier=nav-slam

[Install]
WantedBy=multi-user.target
