# systemd/nav-lidar.service
[Unit]
Description=Navigation - Livox LiDAR Driver
After=nav-lidar-network.service
Requires=nav-lidar-network.service

[Service]
Type=simple
User=sunrise
Group=sunrise
WorkingDirectory=/opt/nav
ExecStart=/bin/bash /opt/nav/scripts/services/nav-lidar.sh
Restart=on-failure
RestartSec=3
TimeoutStopSec=10
KillMode=mixed
KillSignal=SIGINT
SuccessExitStatus=130 SIGINT

Environment=HOME=/home/sunrise
Environment=LANG=en_US.UTF-8

StandardOutput=journal
StandardError=journal
SyslogIdentifier=nav-lidar

[Install]
WantedBy=multi-user.target
