# systemd/nav-autonomy.service
[Unit]
Description=Navigation - Terrain Analysis + Local Planner
After=nav-slam.service
Requires=nav-slam.service

[Service]
Type=simple
User=sunrise
Group=sunrise
WorkingDirectory=/opt/nav
ExecStart=/bin/bash /opt/nav/scripts/services/nav-autonomy.sh
Restart=on-failure
RestartSec=5
TimeoutStopSec=15
KillMode=mixed
KillSignal=SIGINT
SuccessExitStatus=130 SIGINT

Environment=HOME=/home/sunrise
Environment=LANG=en_US.UTF-8

StandardOutput=journal
StandardError=journal
SyslogIdentifier=nav-autonomy

[Install]
WantedBy=multi-user.target
